# scripts/com.lowerrockridge.formspree-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS: cron often cannot read ~/Documents (Operation not permitted). A user LaunchAgent
  runs in your login session and usually works without Full Disk Access for cron.

  Install (edit paths if your home or repo location differs):
    cp scripts/com.lowerrockridge.formspree-sync.plist ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.lowerrockridge.formspree-sync.plist

  Unload:
    launchctl bootout gui/$(id -u) ~/Library/LaunchAgents/com.lowerrockridge.formspree-sync.plist

  Remove the crontab line for the same job to avoid running twice.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.lowerrockridge.formspree-sync</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/USER/Documents/GitHub/rockridge-slow-neighborhood-staging/scripts/sync_supporter_count_formspree.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/USER/Documents/GitHub/rockridge-slow-neighborhood-staging</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>GIT_PUSH</key>
    <string>1</string>
  </dict>
  <key>StartInterval</key>
  <integer>300</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/formspree-sync.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/formspree-sync.log</string>
</dict>
</plist>
